# src/bosixnet-daemon/bosixnet-daemon.service
[Unit]
Description=BOSixNet daemon
After=network.target

[Service]
ExecStart=/usr/bin/bosixnet_daemon

[Install]
WantedBy=multi-user.target
